# profile.d/aliases.sh
alias grep='grep --color=auto'

alias g='git'

alias ls='ls -pG'

alias doma='docker-machine'
alias doco='docker-compose'

alias goodssl='/usr/local/opt/openssl/bin/openssl'
alias libressl='/usr/local/opt/libressl/bin/openssl'

alias vi='vim'

alias an='ansible'
alias anp='ansible-playbook'

alias bse='brew services'
