# profile.d/aliases.sh
alias grep='grep --color=auto'

alias g='git'

alias ls='ls -pG'

alias lockscreen='/System/Library/CoreServices/Menu\ Extras/User.menu/Contents/Resources/CGSession -suspend'

alias gpg='gpg2'

alias doma='docker-machine'
alias doco='docker-compose'
alias bro='google-chrome --incognito $(docker-compose port odoo 8069)'

alias goodssl='/usr/local/opt/openssl/bin/openssl'
alias libressl='/usr/local/opt/libressl/bin/openssl'

alias vi='vim'
