# profile.d/aliases.sh
alias grep='grep --color=auto'

alias g='git'

alias ls='ls -pG'

alias doma='docker-machine'
alias doco='docker-compose'

alias goodssl='/usr/local/opt/openssl/bin/openssl'
alias libressl='/usr/local/opt/libressl/bin/openssl'

alias vi='vim'

alias an='ansible'
alias anp='ansible-playbook'

alias bse='brew services'

alias java7='export JAVA_HOME=$(/usr/libexec/java_home -F -v 1.7); java -version'
alias java8='export JAVA_HOME=$(/usr/libexec/java_home -F -v 1.8); java -version'
alias java9='export JAVA_HOME=$(/usr/libexec/java_home -F -v 9); java -version'
