# profile.d/aliases.sh
alias grep='grep --color=auto'
alias rgrep='grep -r --color=auto'
alias xgrep='grep -r --color=auto --include=\*xml'
alias pgrep='grep -r --color=auto --include=\*py'
alias pogrep='grep -r --color=auto --include=\*po'
alias itpogrep='grep -r --color=auto --include=\*it.po'
alias ygrep='grep -r --color=auto --include=\*yml'
alias csvgrep='grep -r --color=auto --include=\*csv'
alias fegrep='grep -r --color=auto --include=\*feature'
alias agrep='grep -r --color=auto --include=\*xml --include=\*py --include=\*yml --include=\*csv --include=\*mako'

alias pack="ack --type=python"
alias xack="ack --type=xml"

alias b='bzr'
alias git='hub'
alias g='hub'
