# profile.d/exports.sh
export LESS=-FXRSq
export PYTHONSTARTUP=~/.pystartup
export EDITOR=vim
export GPG_TTY=$(tty)
export FZF_DEFAULT_COMMAND=fd
