# profile.d/exports.sh
export LESS=-FXRSq

export PYTHONSTARTUP=~/.pystartup

export HISTCONTROL=erasedups
export HISTSIZE=-1
export HISTFILESIZE=-1

export EDITOR=vim

export GPG_TTY=$(tty)
