# profile.d/exports.sh
export LESS=-FXRSq

export PYENV_ROOT="$HOME/.pyenv"

export HISTCONTROL=erasedups
export HISTSIZE=10000
