# profile.d/path.sh
# on mac, no need to add those. To check again on ubuntu.
# export PATH="$PYENV_ROOT/bin:$HOME/.rbenv/bin:$PATH:/usr/texbin"

if [ -d ~/.linuxbrew/bin ]
then
  export PATH="$HOME/.linuxbrew/bin:$PATH"
fi

export PATH="~/bin:$PATH"
