# profile.d/path.sh
if [ -d ~/Library/Python/2.7/bin ]; then
  export PATH="$PATH:$HOME/Library/Python/2.7/bin"
fi

export PATH="$HOME/bin:/usr/local/sbin:$PATH"
